# build/airootfs/etc/systemd/system/deckos-tts.service
[Unit]
Description=DeckOS Text-to-Speech Service
After=sound.target

[Service]
Type=simple
ExecStart=/usr/local/bin/deckos-tts-server
Restart=on-failure
RestartSec=5
User=deckos
Group=audio
Environment="PIPER_MODEL=/var/lib/deckos/models/piper-en-us"

[Install]
WantedBy=multi-user.target
